# mac-mini/launchd/com.gobag.pluto-rollup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gobag.pluto-rollup</string>

    <!-- Hit the ingester endpoint that runs the rollup. No args = yesterday. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/curl</string>
        <string>-sS</string>
        <string>-X</string>
        <string>POST</string>
        <string>--max-time</string>
        <string>600</string>
        <string>http://127.0.0.1:8765/admin/pluto/rollup</string>
    </array>

    <!-- 03:05 local time every day. Offset 5 min past the hour so we don't
         collide with systemwide wake-ups. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/gobag-pluto-rollup.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/gobag-pluto-rollup.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
